# projects/code.rawkode.academy/forgejo-backup.timer
# Systemd timer for Forgejo backups
[Unit]
Description=Daily Forgejo backup to Cloudflare R2
Requires=forgejo-backup.service

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target