# projects/code.rawkode.academy/forgejo-backup.service
# Systemd service for Forgejo backups
[Unit]
Description=Backup Forgejo to Cloudflare R2
Wants=forgejo.service postgres.service
After=forgejo.service postgres.service

[Service]
Type=oneshot
ExecStart=/opt/forgejo-backup/backup-cronjob.sh
User=root
StandardOutput=journal
StandardError=journal

# Load environment variables from secrets file
EnvironmentFile=/etc/containers/systemd/forgejo-secrets

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp