# projects/code.rawkode.academy/forgejo-backup.service
# Systemd service for Forgejo backups
[Unit]
Description=Backup Forgejo to Cloudflare R2
Wants=forgejo.container
After=forgejo.container

[Service]
Type=oneshot
ExecStart=/home/rawkode/forgejo/backup-cronjob.sh
User=rawkode
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp