# projects/code.rawkode.academy/ngrok.container
[Container]
Image=docker.io/ngrok/ngrok:latest
ContainerName=ngrok
AutoUpdate=registry
LogDriver=journald

Network=forgejo-network
HostName=ngrok

# Mount volumes for configuration
Volume=ngrok-config:/etc/ngrok:Z

# Use secret for auth token
Secret=ngrok-secrets,type=env

# Use the command-line approach instead of config file
Exec=ngrok http 3000 --log stdout

[Unit]
After=network-online.target forgejo-network.service forgejo.service
Requires=forgejo-network.service forgejo.service

[Service]
Restart=always

[Install]
WantedBy=default.target
