# projects/code.rawkode.academy/ngrok.container
[Container]
Image=docker.io/ngrok/ngrok:latest
ContainerName=ngrok
AutoUpdate=registry
LogDriver=journald

Network=forgejo-network
HostName=ngrok

# Mount volumes for configuration
Volume=ngrok-config:/etc/ngrok:Z

# Use secret for auth token
Secret=ngrok-secrets,type=env,target=NGROK_AUTHTOKEN

# Use the command-line approach instead of config file
# Connect to forgejo container on the same network
Exec=ngrok http forgejo:3000 --log stdout

[Unit]
After=network-online.target forgejo.service
Requires=forgejo.service

[Service]
Restart=always

[Install]
WantedBy=default.target
